# production/alert-signs.timer

[Unit]
Description=Backend - timer - alert for missing data

[Timer]
OnBootSec=1min
OnUnitActiveSec=24hours

[Install]
WantedBy=timers.target
