# production/syncRawFluksoDataS.timer

[Unit]
Description=Backend service - timer - test

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
