# production/sync-raw-flukso-data.service
[Unit]
Description=Backend service - sync raw data
After=multi-user.target

[Service]
Type=simple
Restart=no
ExecStart=/opt/vde/venv/bin/python3 "/opt/vde/syncRawFluksoData.py" --mode automatic

[Install]
WantedBy=multi-user.target
