# production/sync-sftp.timer

[Unit]
Description=Backend - timer - sync to sftp server

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
