# production/sync-raw-flukso-data.timer

[Unit]
Description=Backend - timer - sync raw data

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
TimeoutSec=5hours

[Install]
WantedBy=timers.target
