# systemd/vde-daily-recompute.timer
[Unit]
Description=Backend service - daily recompute power data

[Timer]
OnCalendar=*-*-* 1:00:00
Persistent=true

[Install]
WantedBy=timers.target
