# systemd_scripts/sync-raw-flukso-data-local.timer

[Unit]
Description=Backend - timer - sync raw data - test local

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
