# systemd_scripts/syncRawFluksoData.service
[Unit]
Description=Backend service - sync raw data
After=multi-user.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /home/AlexandreHnf/syncRawFluksoData.py --mode automatic

[Install]
WantedBy=multi-user.target