# systemd_scripts/syncRawFluksoDataS.timer
 

[Unit]
Description=Backend service - sync raw data - timer

[Timer]
OnBootSec=5sec
OnUnitActiveSec=3sec

[Install]
WantedBy=timers.target